# periodic_jobs/ai_heartbeat/launchd/ai.erqianwang.agentic-workspace.heartbeat-reflector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.erqianwang.agentic-workspace.heartbeat-reflector</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/erqianwang/Repos/agentic-workspace/periodic_jobs/ai_heartbeat/launchd/run_reflector.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/erqianwang/Repos/agentic-workspace</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/agentic-workspace-heartbeat-reflector.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/agentic-workspace-heartbeat-reflector.err.log</string>
</dict>
</plist>
